# compliance/virtualmachines/roxagent/quadlet/roxagent-watch.path
[Unit]
Description=Watch DNF history for package changes

[Path]
# Rendered at install time from the first known DNF history database path
# present on the target host.
PathModified=@DNF_HISTORY_PATH@
# Trigger the debounce timer, not the service directly.
Unit=roxagent-debounce.timer

[Install]
WantedBy=multi-user.target
