# compliance/virtualmachines/roxagent/quadlet/roxagent-debounce.timer
[Unit]
Description=Debounce timer for reactive roxagent scans

[Timer]
# Fire once, 5 minutes after activation by roxagent-watch.path.
# Multiple DNF operations within the window produce a single scan.
OnActiveSec=5min
# Unload the elapsed timer so future path events can start it again.
RemainAfterElapse=false
# Do not persist across reboots — only react to live changes.
Persistent=false
# Target the reactive container service (generated by Quadlet from
# roxagent-reactive.container).
Unit=roxagent-reactive.service

[Install]
WantedBy=timers.target
